# workshop/admin-tools/upstart_service.conf
description "{{ app_name }}"
  
start on runlevel [3]
stop on shutdown
  
script
    cd /opt/{{ app_name }}/
    java -jar -Daws.user={{ user_name }} {{ app_name }}.jar
end script